# .config/systemd/user/backup-packages.service
[Unit]
Description=Backup list of installed packages

[Service]
ExecStart=%h/bin/backup-packages

