# .config/systemd/user/backup-packages.timer
[Unit]
Description=Daily backup of private data

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
