# .config/systemd/user/backup-packages.service
[Unit]
Description=Backup private thing to mega

[Service]
ExecStart=/home/cyril/bin/backup

