# .config/systemd/user/bum.service
[Unit]
Description=Download covers art played with mopidy

[Service]
ExecStart=/usr/bin/bum --no_display
