# .config/systemd/user/checkmail.service
[Unit]
Description=Run mbsync and notmuch then auto add cal events

[Service]
Type=oneshot
ExecStart=%h/.dotfiles-private/.local/bin/checkmail
