# .config/systemd/user/library-repos.service
[Unit]
Description=Update library repositories

[Service]
Type=oneshot
ExecStart=%h/bin/library-repos

[Install]
WantedBy=multi-user.target
