# .config/systemd/user/pinboard-backup.timer
[Unit]
Description=Dump Pinboard bookmarks timer

[Timer]
Unit=pinboard-backup.service
OnCalendar=daily

[Install]
WantedBy=timers.target
