# .config/systemd/user/restic-check.timer
[Unit]
Description=Restic check

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
