# .config/systemd/user/sway-inactive-window-transparency.service
[Unit]
Description=Make inactive windows in sway semi-transparent
PartOf=sway-session.target
After=sway-session.target

[Service]
ExecStart=%h/.local/bin/sway-inactive-window-transparency
Restart=always
RestartSec=10s

[Install]
WantedBy=sway-session.target
