# .config/systemd/user/tarsnap.timer
[Unit]
Description=Tarsnap timer

[Timer]
Unit=tarsnap.service
OnCalendar=daily

[Install]
WantedBy=timers.target
