# .config/systemd/user/vault.path
[Unit]
Description=Vault

[Path]
PathModified=/home/cyril/Vault

[Install]
WantedBy=multi-user.target
