# .config/systemd/user/waybar-updates.timer
[Unit]
Description=Waybar packages updates notifications async

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
