# .config/systemd/user/waybar-updates.service
[Unit]
Description=Waybar packages updates notifications async

[Service]
ExecStart=/home/cyril/bin/waybar-updates

