# etc/systemd/system/usbguard.service.d/override.conf
[Service]
CapabilityBoundingSet=CAP_DAC_OVERRIDE
ReadWritePaths=-/etc/usbguard/rules.conf
