# .config/systemd/user/poweralertd.service.d/override.conf
[Unit]
PartOf=graphical-session.target
Requires=graphical-session.target dbus.service
After=graphical-session.target dbus.service

[Service]
ExecStart=
ExecStart=/usr/bin/poweralertd