# etc/systemd/system/updatedb.timer.d/updatedb.timer.conf
[Timer]
OnCalendar=
OnCalendar=13:05
Persistent=false
