# etc/systemd/system/man-db.timer.d/man-db.timer.conf
[Timer]
OnCalender=
OnCalendar=13:00
Persistent=false
