# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/usr/bin/agetty --nohostname --nohints --noissue --autologin cyril --noclear %I $TERM
