# .config/systemd/user/wluma.service.d/override.conf
[Unit]
PartOf=sway-session.target lock.target
After=sway-session.target lock.target
StartLimitIntervalSec=300
StartLimitBurst=0

[Service]
Restart=always
RestartSec=5
