# etc/NetworkManager/dispatcher.d/99refresh-py3status
#!/bin/sh
py3-cmd refresh "external_script nmtrust" "wifi" "frame network" "group netrate"
