# etc/systemd/system/privoxy.service.d/override.conf
[Service]
PrivateDevices=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=yes
MountFlags=slave
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
NoNewPrivileges=yes
SystemCallArchitectures=native