# etc/udev/rules.d/83-yubikey.rules
ACTION=="add", SUBSYSTEM=="usb", ENV{ID_VENDOR_ID}=="1050", ENV{ID_MODEL_ID}=="0116", RUN+="/bin/sh -c '/sbin/systemctl stop pcscd.service'"
