# etc/systemd/system/system-dotfiles-sync.service
[Unit]
Description=Sync system dotfiles

[Service]
Type=oneshot
ExecStart=/home/cyril/.dotfiles/setup-system.sh -r
